# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm -D

[program:caddy]
command=/usr/bin/caddy --conf /etc/Caddyfile
